# ops/systemd/kitz-dashboard.service
[Unit]
Description=AI Kitz Labs Dashboard (Next.js production server)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
Type=simple
SyslogIdentifier=kitz-dashboard
UMask=0027

# Template default. Prefer deploying the repo to /opt/kitz-dashboard (or similar)
# instead of user home directories.
WorkingDirectory=/opt/kitz-dashboard

# Secrets belong in an EnvironmentFile, not in the unit file.
# Point this at your deployment-specific env file.
EnvironmentFile=-/etc/kitz-dashboard/kitz-dashboard.env
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=HOME=/var/lib/kitz-dashboard

# Optional: path to a 1Password op:// env template file (non-secret).
Environment=KITZ_OP_ENV_FILE=/etc/kitz-dashboard/kitz-dashboard.op.env
# Optional: 1Password runtime mode (off|auto|required). Default is auto.
Environment=KITZ_1PASSWORD_MODE=auto

ExecStartPre=/usr/bin/test -d /opt/kitz-dashboard/.next
ExecStart=/opt/kitz-dashboard/scripts/start-production.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30

# Run as a dedicated user in template deployments.
User=kitz
Group=kitz

StateDirectory=kitz-dashboard
LogsDirectory=kitz-dashboard
RuntimeDirectory=kitz-dashboard

# Hardening baseline. Adjust ReadWritePaths to match your KITZ_DB_PATH/KITZ_STATE_DIR.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
ReadWritePaths=/opt/kitz-dashboard/.next /var/lib/kitz-dashboard /var/log/kitz-dashboard

[Install]
WantedBy=multi-user.target
